# cloudflare-dyndns.timer
[Unit]
Description=Run cloudflare-dyndns every so often

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
